# folder/user/bin/singb/singb-installer-autoupdater
#!/bin/sh

AUTOSTART_PATH="/etc/init.d/singb-autoupdater"
LOG_PATH="/var/log/singb-autoupdater.log"
PID_FILE="/var/run/singb-autoupdater.pid"
SINGB_UPDATER="/usr/bin/singb/singb-updater"
CONFIG_DIR="/etc/sing-box"
DEFAULT_INTERVAL=3600  # 1 час по умолчанию

install_autoupdater() {
    local interval="${1:-$DEFAULT_INTERVAL}"
    
    cat << EOF > $AUTOSTART_PATH
#!/bin/sh /etc/rc.common
START=99
STOP=10

start() {
    echo "Starting Singb Auto-Updater" > $LOG_PATH
    (while true; do
        for config in config.json config2.json config3.json; do
            url_file="\${CONFIG_DIR}/url_\${config}"
            if [ -f "\$url_file" ]; then
                url=\$(cat "\$url_file")
                if [ -n "\$url" ]; then
                    $SINGB_UPDATER "\$url" "\${CONFIG_DIR}/\$config" >> $LOG_PATH 2>&1
                fi
            fi
        done
        sleep $interval
    done) &
    echo \$! > $PID_FILE
}

stop() {
    [ -f $PID_FILE ] && kill \$(cat $PID_FILE) && rm -f $PID_FILE
    echo "Singb Auto-Updater stopped" >> $LOG_PATH
}

restart() {
    stop
    start
}

status() {
    if [ -f $PID_FILE ] && kill -0 \$(cat $PID_FILE) 2>/dev/null; then
        echo "running"
    else
        echo "stopped"
    fi
}

EOF

    chmod +x $AUTOSTART_PATH
    $AUTOSTART_PATH enable
    $AUTOSTART_PATH start
}

uninstall_autoupdater() {
    $AUTOSTART_PATH stop
    $AUTOSTART_PATH disable
    rm -f $AUTOSTART_PATH $PID_FILE
    echo "Singb Auto-Updater uninstalled" > $LOG_PATH
}

case "$1" in
    install)
        install_autoupdater "$2"
        ;;
    uninstall)
        uninstall_autoupdater
        ;;
    status)
        $AUTOSTART_PATH status
        ;;
    *)
        echo "Usage: $0 {install [interval in seconds] | uninstall | status}"
        ;;
esac
